# cron/comments_scheduler.cron
#!/bin/bash
source ve/bin/activate && python ummeli/manage.py comments_scheduler --settings=vlive_settings && deactivate
