# cron/ummeli.jobs.cron
#!/bin/bash
source ve/bin/activate && python ummeli/manage.py ummeli_jobs_cron --settings=vlive_settings && deactivate
